# etc/supervisor/kotti.conf
[program:kotti]
autorestart=true
user=kotti
command=/home/kotti/maptheuk/bin/uwsgi --ini-paste /home/kotti/maptheuk/app.ini
directory=/home/kotti/maptheuk
redirect_stderr=true
